# profile.d/profile_SunOS
if [ $TERM != "dumb" ]; then
USERID=`id |cut -f2 -d "="|cut -f1 -d "("`
case "$0" in
-sh | -ksh | -jsh | -bash)
        if [ -z "$USER" ] ; then
           USER=$LOGNAME
        fi
        if [ -z "$LANG" ] ; then
           LANG=C
        fi
        echo "Current User:"$USER "($LANG)"
        if [ $USERID = "0" ] ; then
		LANG=C  echo "正在以root超级用户运行，此用户执行的动作可能使系统造成无法恢复的损坏!!"
        fi
esac
  
fi
